# naga-remap.service
[Unit]
Description=Razer Naga V2 Side Button Remapper
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/bin/naga-remap
Restart=on-failure
RestartSec=5

# Hardening
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
